# src/workshop/profile.d/profile-environment-variables.sh
export GEM_HOME=/usr/local/bundle
export PATH=$PATH:$PWD/tools/bin:/usr/local/bundle/bin