# src/devenv/profile.d/profile-environment-variables.sh
export GEM_HOME=/usr/local/bundle
export AWS_DEFAULT_REGION="eu-central-1"
export AWS_DEFAULT_PROFILE="biotronik-poc-dev"
export INFRA_ENVIRONMENT="dev"
export PATH=$PATH:$PWD/tools/bin:/usr/local/bundle/bin